# rpi_zero_image_build_pipeline/sa-image/device/rootfs-overlay/etc/systemd/system/pi-controller.service
[Unit]
Description=Pi Controller Node.js Service
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/pi-controller.env
User=controller
Group=controller
WorkingDirectory=/opt/pi-controller-nodejs
ExecStart=/usr/local/bin/node /opt/pi-controller-nodejs/dist/app.js --trace-warnings --expose-gc
Restart=always
RestartSec=5
# Give a bit of memory headroom on Zero
OOMScoreAdjust=-900
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
